# capercio/config/autostart/kiosk.desktop
[Desktop Entry]
Type=Application
Name=Kiosk
Comment=
Exec=/home/pi/bin/kiosk.sh
Terminal=false
Hidden=false
